# roles/custom/nextcloud_server/templates/systemd/nextcloud-cron.timer.j2
[Unit]
Description=Launches nextcloud-cron.service periodically

[Timer]
Unit=nextcloud-cron.service
OnCalendar={{ nextcloud_nextcloud_cron_schedule }}

[Install]
WantedBy=timers.target
