# roles/custom/nextcloud_nginx_proxy/templates/systemd/nextcloud-nginx-proxy-ssl-certificate-renew.service.j2
[Unit]
Description=Renews Let's Encrypt SSL certificates periodically

[Service]
Type=oneshot
ExecStart={{ nextcloud_nginx_proxy_bin_path }}/ssl-certificates-renew
