# roles/nextcloud-server/templates/systemd/nextcloud-preview-generator.timer.j2
[Unit]
Description=Run the generate-preview app of Nextcloud (timer)

[Timer]
{{ nextcloud_previewgenerator_timer }}
Persistent=true

[Install]
WantedBy=timers.target
