# roles/nextcloud-server/templates/cron.d/nextcloud-cron.j2
MAILTO="{{ ssl_support_email }}"

*/15 * * * * root /usr/bin/docker exec nextcloud-apache su - www-data -s /bin/bash -c 'php /var/www/html/cron.php'
