# roles/nextcloud_server/templates/cron.d/nextcloud-cron.j2
MAILTO="{{ nextcloud_ssl_support_email }}"

*/15 * * * * root /usr/bin/docker exec -u www-data nextcloud-apache php /var/www/html/cron.php
